# sources/meta-development/recipes-multimedia/bluealsa/files/bluealsa.service
[Unit]
Description=Start the bluealsa service
After=bluetooth.service
After=systemd-remount-fs.service

[Service]
ExecStart=/usr/bin/bluealsa

[Install]
WantedBy=bluetooth.target
